# src/rootfs/overlay/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
#all mounts should be done in initamfs
null::sysinit:/bin/hostname -F /etc/hostname
# now run any rc scripts
[% IF ENV.UPD_CONFIG == "rel" -%]
::sysinit:/etc/init.d/rcS start
[% ELSE -%]
::sysinit:time /etc/init.d/rcS start
::sysinit:/bin/echo "Rootfs loading time (/etc/init.d/S*)." > /dev/console
[% END -%]

# Set up a couple of getty's
#tty1::respawn:/sbin/getty 38400 tty1
#tty2::respawn:/sbin/getty 38400 tty2

# Put a getty on the serial port
#ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100 # GENERIC_SERIAL

# Logging junk
#tty3::respawn:/usr/bin/tail -f /var/log/messages

[% IF CONFIG_TESTSERVER_ENABLE || CONFIG_TESTTOOL_ENABLE -%]
::sysinit:/opt/elecard/bin/test.sh start
[% END -%]

[% IF ! ENV.ROOTFS_DISABLE_CONSOLE_SHELL -%]
[% IF ENV.ROOTFS_DISABLE_AUTO_LOGIN -%]
::respawn:-/bin/login </dev/console
[% ELSE -%]
::respawn:-/bin/sh --login
[% END -%]
[% END -%]

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
null::shutdown:/usr/bin/killall klogd
null::shutdown:/usr/bin/killall syslogd
null::shutdown:/bin/umount /var/etc
null::shutdown:/bin/sleep 1
null::shutdown:/bin/umount -a -r
null::shutdown:/sbin/swapoff -a

